# etc/profile.d/wazo-apt-aliases.sh
alias apt-get='/usr/bin/apt-get-warning'
alias apt='/usr/bin/apt-warning'
